# swbatch/swbatch.conf
[program:swbatch]
command=/var/swbatch/swbatch
autostart=true
autorestart=true
stderr_logfile=/var/log/swbatch.err.log
stdout_logfile=/var/log/swbatch.out.log

